# etc/systemd/system/steamos-startup.service
[Unit]
Description=Set autologin-session to gamescope-session
After=local-fs.target
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/bin/steamos-startup-gamescope

[Install]
WantedBy=multi-user.target
